# files/etc/profile.d/cachyos-boppos-wayland.sh
# Forces Electron apps (Discord, VS Code, Obsidian) to run natively on Wayland.
export ELECTRON_OZONE_PLATFORM_HINT="auto"